# backend/deploy/pixel-social-world.service
[Unit]
Description=Pixel Social World backend
Wants=network-online.target
After=network-online.target postgresql.service redis-server.service

[Service]
Type=simple
User=pixelsocial
Group=pixelsocial
WorkingDirectory=/opt/pixel-social-world/backend
EnvironmentFile=/etc/pixel-social-world/backend.env
ExecStartPre=/opt/pixel-social-world/backend/bin/pixel-social-world-preflight -strict
ExecStart=/opt/pixel-social-world/backend/bin/pixel-social-world-server
Restart=always
RestartSec=3
KillSignal=SIGTERM
TimeoutStopSec=15
LimitNOFILE=1048576
NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=full
ProtectHome=true
ReadWritePaths=/opt/pixel-social-world /var/lib/pixel-social-world

[Install]
WantedBy=multi-user.target
