# backend/deploy/pixel-social-world-retention-cleanup.timer
[Unit]
Description=Run Pixel Social World retention cleanup daily

[Timer]
OnCalendar=*-*-* 04:15:00
Persistent=true
RandomizedDelaySec=15m

[Install]
WantedBy=timers.target
